# systemd/.config/systemd/user/himalaya-watcher.service
[Unit]
Description=Himalaya watcher

[Service]
ExecStart=himalaya watch
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
