# systemd/.config/systemd/user/mail.timer
[Unit] 
Description=Syncing mailbox

[Timer] 
OnCalendar=*:0/5
Persistent=true
OnBootSec=1 m

[Install] 
WantedBy=timers.target
